# supervisord.conf
[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface
[supervisorctl]
serverurl=unix:///run/supervisord.sock ; use a unix:// URL  for a unix socket
[program:redis-4.0.11]
  command=/bin/sh -c "exec /usr/local/redis/bin/redis-server /etc/redis/redis.conf"
  autostart=true
  autorestart=true
  startretries=3
  startsecs=10 
  stdout_events_enabled=true
  stderr_events_enabled=true
  stderr_logfile_maxbytes=20MB
  stdout_logfile_backups = 5
